# php7.2/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm -D
autostart=true
autorestart=true
priority=5